# scripts/el6/cif_crontool.cron
# set the path
PATH=/bin:/usr/bin:/usr/local/bin:/opt/cif/bin

# pull feed data
05     *       * * * 	cif	/opt/cif/bin/cif_crontool -p hourly -P -d -N -A root >> /var/log/cif/crontool_hourly.log 2>&1
30     00      * * * 	cif	/opt/cif/bin/cif_crontool -p daily -P -d -N -A root >> /var/log/cif/crontool_daily.log 2>&1

# if you've enabled feed generation in your config
45     *       * * * 	cif	/opt/cif/bin/cif_feed -d >> /var/log/cif/cif_feed.log 2>&1
